# iniParser/sample.ini
[program:proc1]
command=/app/bin/start.sh
stdout_logfile=/var/log/proc1_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/proc1_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
autostart=true
autorestart=true
priority=1
stopsignal=QUIT
stdout_events_enabled=true
stderr_events_enabled=true

[program:proc2]
command=/app/bin/proc2 -config /app/conf/test.conf
stdout_logfile=/var/log/proc2_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/proc2_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
autostart=true
autorestart=true
priority=2
stdout_events_enabled=true
stderr_events_enabled=true

[program:proc3]
command=/app/bin/proc3 -config /app/conf/test.conf
stdout_logfile=/var/log/proc3_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/proc3_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
autostart=true
autorestart=true
priority=3
stdout_events_enabled=true
stderr_events_enabled=true